# infrastructure/ansible/roles/pbx-host/templates/pbx-healthcheck.timer.j2
[Unit]
Description=PBX Platform periodic health checks

[Timer]
OnCalendar={{ health_check_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
