# infrastructure/ansible/roles/pbx-host/templates/pbx-backup.timer.j2
[Unit]
Description=PBX Platform daily backup timer

[Timer]
OnCalendar={{ backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
